# configs/client-vnc/rootfs/etc/services.d/xvfb/run
#!/usr/bin/with-contenv sh

exec /usr/bin/Xvfb $DISPLAY -screen 0 ${DISPLAY_WIDTH}x${DISPLAY_HEIGHT}x24
